# sappho-pda/Tools/sapphocape.service
[Unit]
Description=SAPPHO-PDA Cape Service

[Service]
ExecStart=/home/debian/sappho-pda/Tools/capeMod.sh

[Install]
WantedBy=multi-user.target
